# estkme-bark-monitor.service
[Unit]
Description=ESTKme Bark stock monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/estkme-bark-monitor
EnvironmentFile=/opt/estkme-bark-monitor/.env
ExecStart=/usr/bin/python3 /opt/estkme-bark-monitor/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
